# system/etc/systemd/system/user@.service.d/99-oomd.conf
# Kill user-session processes when memory pressure stays high.
# Triggered when >50% pressure persists for >20s on the user slice.
[Service]
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=50%
